# backend/garage-controller.service
[Unit]
Description=Garage Door Controller (gunicorn + uvicorn)
After=network.target

[Service]
User=pi
Group=www-data
WorkingDirectory=/home/pi/garageController
Environment="PATH=/home/pi/venv/bin"
ExecStart=/home/pi/venv/bin/gunicorn main:app -w 1 -k uvicorn.workers.UvicornWorker --bind 0.0.0.0:8080 --timeout 120 --graceful-timeout 30
Restart=always
RestartSec=5
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
